# systemd/.config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper Setter

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/home/balaji/bin/init-wallpaper
